# agent-backend.service
[Unit]
Description=Enterprise Knowledge Base Q&A Agent Backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/agent_rag_project
Environment="PATH=/opt/agent_rag_project/venv/bin"
Environment="ZHIPU_API_KEY=your_zhipu_api_key_here"

# 启动 FastAPI 服务
ExecStart=/opt/agent_rag_project/venv/bin/uvicorn backend.main:app --host 0.0.0.0 --port 8000 --workers 2

# 重启策略
Restart=always
RestartSec=10

# 安全配置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
